# deploy/scalping.service
[Unit]
Description=Scalping Radar - Outil d'aide a la decision
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/scalping
ExecStartPre=-/usr/bin/docker rm -f scalping-radar
ExecStartPre=-/usr/bin/mkdir -p /opt/scalping/data
ExecStart=/usr/bin/docker run --name scalping-radar --env-file /opt/scalping/.env -v /opt/scalping/data:/app/data -v /opt/scalping/v2-assets:/app/frontend-react/dist/assets:ro -p 127.0.0.1:8000:8000 scalping-radar:latest
ExecStop=/usr/bin/docker stop scalping-radar
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
